# wallgarden/service/wallgarden.service.template
[Unit]
Description=Wallgarden slideshow

[Service]
ExecStart=__sys.executable__ __wallgarden.cli__ random --pinned

[Install]
WantedBy=default.target
