# wallgarden/service/wallgarden.timer.template
[Unit]
Description=Wallgarden slideshow timer

[Timer]
OnBootSec=1min
OnUnitActiveSec=__minutes__min

[Install]
WantedBy=timers.target
